# updatecli/new_stuff/updatecli.d/services/check_github_releases.service
[Unit]
Description=Check GitHub releases
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/updatecli apply --config /home/jamtools/code/updatecli_github_commit.yml
EnvironmentFile=/home/jamtools/code/secrets.env

[Install]
WantedBy=multi-user.target
