# updatecli/new_stuff/updatecli.d/services/check_github_releases.timer
[Unit]
Description=Timer for checking GitHub releases

[Timer]
OnBootSec=1min
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
